# gekko/gekko.conf
description "gekko"
author "jonny64"
start on filesystem and started docker.io
stop on runlevel [!2345]
respawn
script
  sudo /usr/bin/docker start -a gekko
end script
